# App/PrivilegedHelper/container.matrixreligio.com.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>container.matrixreligio.com.helper</string>
    <key>MachServices</key>
    <dict>
        <key>container.matrixreligio.com.helper</key>
        <true/>
    </dict>
    <key>BundleProgram</key>
    <string>Contents/Library/PrivilegedHelperTools/container.matrixreligio.com.helper</string>
    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
